# deploy/pantry-display.service
[Unit]
Description=Pantry Display Flask App
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/pantry-display
# venv first, then system tools (waveshare epdconfig uses cat/grep to detect Pi)
Environment=PATH=/home/pi/pantry-display/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/pi/pantry-display/.venv/bin/python app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
